# etc/init.d/fprintd-openrc
#!/sbin/openrc-run

description="Fingerprint Authentication Daemon"
command="/usr/libexec/fprintd"
command_args=""
command_background="yes"
pidfile="/run/fprintd.pid"

depend() {
    need dbus
    after dbus
}

start_pre() {
    # Ensure the state directory exists with the correct permissions
    checkpath --directory --owner root:root --mode 0700 /var/lib/fprint

    # Check and create necessary device permissions
    # These may need to be set up manually or via udev rules depending on your system
}

start() {
    ebegin "Starting fprintd"
    # Start the daemon with confinement if possible
    start-stop-daemon --start --exec ${command} -- ${command_args}
}
