# etc/init.d/fprintd
#!/sbin/openrc-run
# Copyright (c) 2007-2021 The Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Fingerprint Authentication Daemon"
command="/usr/libexec/fprintd"
pidfile="/run/fprintd.pid"

depend() {
	# fprintd is a D-Bus service and requires it to be running.
	need dbus
}

start_pre() {
	# Ensure the state directory exists with secure permissions.
	checkpath --directory --owner root:root --mode 0700 /var/lib/fprint
}
